# supervisord/container_config_folder/supervisord-sshd.conf
[program:sshd]
command=/scripts/start_sshd.sh
autostart=yes
autorestart=true
sterr_logfile=/var/log/supervisord/sshd-err.log
stout-logfile=/var/log/supervisord/sshd.log


