# xfce4/xfce4-notworking/etc/supervidor.d/vncserver.ini
#Based on: https://github.com/jumanjiman/docker-xrdp/blob/master/fedora/etc/supervisord.d/xrdp.ini
[program:vncserver]
command=/usr/bin/vncserver -fg
process_name = vncserver
